# pika-gnome-settings/etc/xdg/autostart/pika-theme.desktop
[Desktop Entry]
Name=Pika Theme Setup
Exec=pika-theme-setup
Terminal=false
Type=Application
