# pika-settings/etc/xdg/autostart/pika-theme.desktop
[Desktop Entry]
Type=Application
Name=Pika Theme Setup
Exec=pika-theme-setup
